# agent/packaging/linux/postinstall.sh
#!/usr/bin/env sh
set -eu

echo "Run 'ov-agent install' to bind this device, then enable the user service with:"
echo "systemctl --user enable --now ov-computeruse-agent.service"
